# ops/systemd/docflow-backup.timer
[Unit]
Description=Run DocFlow backup on schedule

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true
Unit=docflow-backup.service

[Install]
WantedBy=timers.target
